# docker/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/var/run/supervisord.pid
minfds=1024
minprocs=200

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[program:mysqld]
command=/usr/bin/mysqld_safe

[program:postfix]
command = /usr/sbin/postfix start
startsecs = 0
autorestart = false

[program:memcached]
command=/usr/bin/memcached -u memcached
stderr_logfile=/var/log/supervisor/memcached.log
stdout_logfile=/var/log/supervisor/memcached.log
